# artcake.service
[Unit]
Description=ArtCake Web Application
After=network.target

[Service]
Type=simple
User=artcake
WorkingDirectory=/opt/artcake
ExecStart=/usr/bin/java -jar /opt/artcake/ArtCake-0.0.1-SNAPSHOT.jar --spring.profiles.active=prod
Restart=always
RestartSec=10

# Environment variables (adjust as needed)
Environment=DATABASE_URL=jdbc:postgresql://ider-database.westeurope.cloudapp.azure.com:5433/h184940
Environment=DATABASE_USERNAME=h184940
Environment=DATABASE_PASSWORD=your_db_password
Environment=MAIL_PASSWORD=your_mail_password
Environment=APP_BASE_URL=https://www.artcake.no
Environment=MAIL_HOST=mail.artcake.no
Environment=MAIL_PORT=465
Environment=MAIL_USERNAME=artcake@artcake.no
Environment=APP_MAIL_FROM=artcake@artcake.no

# Java options for production
Environment=JAVA_OPTS=-Xmx512m -Xms256m -server

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=artcake

[Install]
WantedBy=multi-user.target
